# debian/icaas-agent.icaas.upstart
# icaas - Image Creator as a Service

description "ICaaS agent"

start on runlevel [2345]
stop on runlevel [!2345]

umask 022

env ICAAS_MONITOR_SIGSTOP=1
expect stop

pre-start script
    test -f /.icaas || { stop; exit 0; }
    test -f /etc/icaas/manifest.cfg || { stop; echo "The manifest file missing!" >&2 ; exit 1; }

    # Wait for the DNS to come up
    for i in 0 1 2 4 8; do
        sleep $i
        if /usr/bin/curl -4 -I -s -S http://www.google.com > /dev/null; then
             exit 0
        fi
    done
    stop
    echo "Network connection failed!" >&2
    exit 2
end script

exec /usr/bin/icaas -m /etc/icaas/manifest.cfg

post-start script
    rm -f /.icaas
end script

# vim: set sta sts=4 shiftwidth=4 sw=4 et ai syntax=upstart:
